# 5-victorialogs/victoria-logs.service
[Unit]
Description=VictoriaLogs
After=network.target

[Service]
Type=simple
User=victorialogs
ExecStart=/usr/local/bin/victoria-logs \
  -storageDataPath=/var/lib/victoria-logs/data \
  -retentionPeriod=30d \
  -httpListenAddr=:9428 \
  -logNewStreams \
  -logIngestedRows
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target